# packages/sssd/debian/sssd-common.prerm
#! /bin/sh -e

invoke_failure() {
    # invoke-rc.d failed
    return
}

#DEBHELPER#

