# packages/autopoweroff/debian/postinst
#!/bin/bash

update-rc.d autopoweroff defaults

cat >/usr/bin/autopoweroff <<EOS
#!/bin/bash
gksu /usr/sbin/autopoweroff
EOS

chmod 755 /usr/bin/autopoweroff
