# packages/default/puavo-load-reporter/debian/puavo-load-reporter.upstart
description "puavo-load-reporter"

start on runlevel [23]
stop on shutdown

exec /usr/sbin/puavo-load-reporter

respawn
