# packages/puavo-users/debian/puavo-rest-bootserver.upstart
description "puavo-rest-bootserver"

start on runlevel [23]
stop on shutdown

pre-start script
    mkdir -p /run/puavo-rest-bootserver
    chown puavo:puavo /run/puavo-rest-bootserver
end script

chdir /var/app/puavo-rest
env RACK_ENV=production

exec sudo -E -u puavo bundle exec rackup -o 127.0.0.1

respawn
